# watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
   "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.atakan.arxivwatcher</string>

  <!-- Run every day at 9:30 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/takmete/arxiv_watcher.py</string>
    <string>--hours</string>
    <string>24</string>
    <string>--out</string>
    <string>/Users/takmete/Papers/arxiv_hits</string>
    <string>--notify</string>
  </array>

  <!-- Optional: logs -->
  <key>StandardOutPath</key>
  <string>/Users/takmete/Library/Logs/arxiv_watcher.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/takmete/Library/Logs/arxiv_watcher.err</string>
</dict>
</plist>
